# app/templates/production/_upstart.conf
# Manager {{ lname|capitalize }} upstart script
# To be saved as /etc/init/manager-{{ lname }}.conf
#

description     "Manager {{ lname|capitalize }}"
author          "Developer <developer@muzzley.com>"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [06]

# Prepare the environment
#   Create directories for logging and process management
#   Change ownership to the user running the process
pre-start script
    mkdir -p {logPath}/[packageName]
    chown -R node:node {logPath}
end script

# If the process quits unexpectedly trigger a respawn
respawn

# Start the process
exec su - node -c 'nvm use 0.10 ; cd {currentPath} && umask 0002 && node {currentPath}/index.js >> {logPath}/[packageName]/[packageName].log 2>&1'
